# srcs/ftps/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:setup]
command= vsftpd /etc/vsftpd/vsftpd.conf

[program:telegraf]
command=/telegraf/usr/bin/telegraf --config /etc/telegraf/telegraf.conf