# srcs/mysql/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:setup]
command= sh init_docker.sh

[program:telegraf]
command=/telegraf/usr/bin/telegraf --config /etc/telegraf/telegraf.conf