# srcs/grafana/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:setup]
command=/grafana-7.0.3/bin/grafana-server

[program:telegraf]
command=/telegraf/usr/bin/telegraf --config /etc/telegraf/telegraf.conf